# build/docker/supervisor.ini
[supervisord]
user=root
pidfile=/var/run/supervisord.pid
stdout_logfile=/proc/self/fd/2
stdout_logfile_maxbytes=0
stderr_logfile=/proc/self/fd/2
stderr_logfile_maxbytes=0

[program:install]
command=/bin/bash -c 'if [ -z "$( ls -A /opt/pseudify/userdata/ )" ]; then rsync -raq /opt/pseudify/userdata_template/ /opt/pseudify/userdata/; else rm -rf /opt/pseudify/userdata/var/; fi && composer install --no-dev --optimize-autoloader --classmap-authoritative --no-progress --no-interaction --quiet && bin/pseudify importmap:install && bin/pseudify asset-map:compile'
directory=/opt/pseudify/core
priority=1
autorestart=false
startretries=0
user=pseudify
stdout_logfile=/proc/self/fd/2
stdout_logfile_maxbytes=0
stderr_logfile=/proc/self/fd/2
stderr_logfile_maxbytes=0

[program:php-fpm]
command=php-fpm -F
priority=2
autorestart=true
user=root
stdout_logfile=/proc/self/fd/2
stdout_logfile_maxbytes=0
stderr_logfile=/proc/self/fd/2
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g "daemon off;"
priority=2
autorestart=true
user=root
stdout_logfile=/proc/self/fd/2
stdout_logfile_maxbytes=0
stderr_logfile=/proc/self/fd/2
stderr_logfile_maxbytes=0
